# sparse/lib/systemd/system/droid-cameradetect.service
[Unit]
Description=Detect camera module change
After=init-done.service
Conflicts=shutdown.target actdead.target

[Service]
Type=oneshot
ExecStart=/bin/sh /usr/bin/droid/droid-cameradetect.sh
RemainAfterExit=yes
